# docker/pre/supervisord.conf
[supervisord]
nodaemon=true

[program:octane]
process_name=%(program_name)s
command=php /var/www/html/artisan octane:start --port=9000 --host=0.0.0.0
autostart=true
autorestart=true
